# packages/linux-api/systemd/things-bridge-api.service
# This is a reference template. The install wizard (bun index.ts) generates this automatically with real values.
[Unit]
Description=Things Bridge API Service
After=network.target

[Service]
Type=simple
User=things-bridge
WorkingDirectory=/path/to/things_bridge/packages/linux-api
ExecStart=/usr/local/bin/bun /path/to/things_bridge/packages/linux-api/src/index.ts
Restart=always
RestartSec=10

[Install]
WantedBy=multi-user.target
