# packages/mac-agent/launchd/com.things-bridge.agent.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<!-- This is a reference template. The install wizard (bun index.ts) generates this automatically with real values. -->
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.things-bridge.agent</string>
    <key>ProgramArguments</key>
    <array>
        <string>/opt/homebrew/bin/bun</string>
        <string>/path/to/things_bridge/packages/mac-agent/src/index.ts</string>
    </array>
    <key>WorkingDirectory</key>
    <string>/path/to/things_bridge/packages/mac-agent</string>
    <key>EnvironmentVariables</key>
    <dict>
        <key>PATH</key>
        <string>/opt/homebrew/bin:/usr/local/bin:/usr/bin:/bin</string>
    </dict>
    <key>RunAtLoad</key>
    <true/>
    <key>KeepAlive</key>
    <true/>
    <key>StandardOutPath</key>
    <string>/tmp/things-bridge-agent.log</string>
    <key>StandardErrorPath</key>
    <string>/tmp/things-bridge-agent.err</string>
</dict>
</plist>
